# supervise/main.conf
[program:main]
command=python run.py
directory=/prosecco
