# supervise/player.conf
[program:master_player]
command=python master_player.py
directory=/prosecco
priority=500
